# roles/dcos-filebeat/templates/systemd/logrotate.service.j2
[Unit]
Description=Rotate and Compress System Logs

[Service]
ExecStart=/usr/sbin/logrotate {{ logrotate_conf_file }}
