# roles/dcos-filebeat/templates/systemd/logrotate.timer.j2
[Unit]
Description=Docker App Log Rotation

[Timer]
OnCalendar=30m
AccuracySec=1m
Persistent=true
